# system/etc/systemd/user/battery-check.service
[Unit]
Description=Battery Check Service
BindsTo=desktop.target

[Service]
ExecStart=battery-check
