# system/etc/systemd/user/batter-check.service
[Unit]
Description=Battery Check Service

[Service]
ExecStart=battery-check

[Install]
RequiredBy=desktop.target
