# system/etc/systemd/user/battery-check.service
[Unit]
Description=Battery Check Service
Wants=desktop.timer

[Service]
Type=oneshot
ExecStart=/usr/bin/battery-check
