# system/etc/systemd/user/dunst.service
[Unit]
Description=Dunst Service
BindsTo=desktop.target

[Service]
ExecStart=/usr/bin/dunst
