# system/etc/systemd/user/dunst.service
[Unit]
Description=Dunst Service

[Service]
ExecStart=/usr/bin/dunst

[Install]
WantedBy=desktop.target
