# system/etc/systemd/user/keyboard.target
[Unit]
Description=Keyboard Target
After=default.target
Requires=touchcursor.target keyboard-layout.service

[Install]
RequiredBy=desktop.target
