# system/etc/systemd/user/nm-applet.service
[Unit]
Description=Network Manager Applet Service

[Service]
ExecStart=/usr/bin/nm-applet

[Install]
RequiredBy=desktop.target
