# system/etc/systemd/user/nm-applet.service
[Unit]
Description=Network Manager Applet Service
BindsTo=desktop.target

[Service]
ExecStart=/usr/bin/nm-applet
Environment='DISPLAY=:0'
