# system/etc/systemd/user/picom.service
[Unit]
Description=Picom Service
BindsTo=desktop.target

[Service]
ExecStart=/usr/bin/picom --config %h/.config/picom.conf
Environment='DISPLAY=:0'
