# system/etc/systemd/user/picom.service
[Unit]
Description=Picom Service

[Service]
ExecStart=/usr/bin/picom --config %h/.config/picom.conf

[Install]
RequiredBy=desktop.target
