# system/etc/systemd/user/touchcursor.target
[Unit]
Description=Touch Cursor Service
BindsTo=keyboard.target
