# system/etc/systemd/user/touchcursor.target
[Unit]
Description=Touch Cursor Service
After=default.target

[Install]
RequiredBy=keyboard.target
