# system/etc/systemd/user/wallpaper.service
[Unit]
Description=Wallpaper Service
BindsTo=desktop.target

[Service]
ExecStart=/usr/bin/nitrogen --restore
