# system/etc/systemd/user/wallpaper.service
[Unit]
Description=Wallpaper Service

[Service]
Type=oneshot
ExecStart=/usr/bin/nitrogen --restore

[Install]
RequiredBy=desktop.target
