# system/etc/systemd/user/xray.service
[Unit]
Description=Xray Service
BindsTo=desktop.target

[Service]
ExecStart=/usr/bin/xray -c %h/vpn/ir4.json
