# system/etc/systemd/system/vpn.target
[Unit]
Description=VPN Target
After=network.target
Requires=xray.service tun2socks.service tuntap-route.service
