# system/etc/systemd/system/vpn.target
[Unit]
Description=VPN Target
After=network.target
Requires=xray.service tun2socks.service

[Install]
WantedBy=multi-user.target
