# system/etc/systemd/system/xray.service
[Unit]
Description=Xray Service

[Service]
ExecStart=/usr/bin/xray -c /etc/vpn/ir4.json
